# usb/udev/97-batu.rules
SUBSYSTEM=="usb", RUN+="/etc/batu/usb/usb.sh"
